# build/deb/reqguard_0.3.20_all/lib/systemd/system/reqguard-firewall.service
[Unit]
Description=Reqguard firewall ban list sync
Wants=network-online.target
After=network-online.target ufw.service nftables.service
StartLimitIntervalSec=2min
StartLimitBurst=12

[Service]
Type=oneshot
EnvironmentFile=-/etc/default/reqguard
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/bin/reqguard sync-firewall
RemainAfterExit=yes
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
