# build/deb/reqguard_0.3.19_all/lib/systemd/system/reqguard-firewall.service
[Unit]
Description=Reqguard firewall ban list sync
After=ufw.service nftables.service network-pre.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/default/reqguard
ExecStart=/usr/bin/reqguard sync-firewall
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
